# deploy/intelhq-api.service
[Unit]
Description=IntelHQ API Server
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/intelhq
ExecStart=/usr/bin/node src-tauri/sidecar/local-api-server.mjs
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=intelhq-api

# Environment
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=API_ROOT=/opt/intelhq/api
EnvironmentFile=-/opt/intelhq/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/intelhq

[Install]
WantedBy=multi-user.target
